# systemd/zombie-monitor.timer
[Unit]
Description=Daily Zombie Process Check

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
